# scripts/systemd/juara-wifi-watchdog.timer
[Unit]
Description=Run Juara WiFi reconnect watchdog

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true
Unit=juara-wifi-watchdog.service

[Install]
WantedBy=timers.target
