# systemd/juara-service-hourly-reboot.service
[Unit]
Description=Juara service-mode hourly reboot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/juara-planned-reboot
