# systemd/juara-daily-reboot.timer
[Unit]
Description=Reboot the Juara station at planned daily checkpoints

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=juara-daily-reboot.service

[Install]
WantedBy=timers.target
